# backend.service
[Unit]
Description=Architecture Academics Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/Suresh_Sir_Arch/backend
Environment="PATH=/home/ec2-user/Suresh_Sir_Arch/backend/venv/bin"
ExecStart=/home/ec2-user/Suresh_Sir_Arch/backend/venv/bin/python run_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
