# deployment/backend.service
[Unit]
Description=Architecture Academics Portal Backend API
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/architecture-academics-portal/backend
Environment=PYTHONPATH=/home/ec2-user/architecture-academics-portal/backend
Environment=HOST=0.0.0.0
Environment=PORT=8000
ExecStart=/home/ec2-user/architecture-academics-portal/backend/venv/bin/python main.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=backend

[Install]
WantedBy=multi-user.target